# com.hanz.discordpresence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hanz.discordpresence</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/janpalenik/Developer/RPC/.venv/bin/python3</string>
        <string>-u</string>
        <string>/Users/janpalenik/Developer/RPC/presence.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/janpalenik/Developer/RPC/presence.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/janpalenik/Developer/RPC/presence.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/janpalenik/Developer/RPC</string>
</dict>
</plist>
